# templates/wyoming-satellite.service.j2
[Unit]
Description=Wyoming Satellite
After=network.target
Wants=network.target

[Service]
Type=simple
User=wyoming
Group=wyoming
WorkingDirectory=/opt/wyoming-satellite
Environment=PATH=/opt/wyoming-satellite/venv/bin:/usr/bin:/bin
ExecStart=/opt/wyoming-satellite/venv/bin/python -m wyoming_satellite \
    --name "{{ wyoming_satellite_name | default(ansible_hostname) }}" \
    --uri "{{ wyoming_satellite_uri | default('tcp://0.0.0.0:10700') }}" \
    --mic-command "{{ wyoming_satellite_mic_command | default('arecord -D plughw:CARD=1,DEV=0 -r 16000 -c 1 -f S16_LE -t raw') }}" \
    --snd-command "{{ wyoming_satellite_snd_command | default('aplay -D plughw:CARD=1,DEV=0 -r 22050 -c 1 -f S16_LE -t raw') }}" \
    --mic-auto-gain 5 \
    --mic-noise-suppression 2
    {% if wyoming_satellite_wake_command is defined %}--wake-command "{{ wyoming_satellite_wake_command }}" \{% endif %}
    {% if wyoming_satellite_wake_uri is defined %}--wake-uri "{{ wyoming_satellite_wake_uri }}" \{% endif %}
    {% if wyoming_satellite_wake_word_name is defined %}--wake-word-name "{{ wyoming_satellite_wake_word_name }}" \{% endif %}
    {% if wyoming_satellite_debug is defined and wyoming_satellite_debug %}--debug{% endif %}

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target